# etc/cron/prod.cron
# Crontab for RWizFlowy
# chkcrontab: disable-msg=USER_NOT_FOUND
APP_HOME=/opt/rwizflowy/prod

# Run daily scripts (00:01)
01 00 * * * $APP_HOME/etc/cron/r_daily_runner.sh

# Run hourly scripts (0*:00)
00 * * * * $APP_HOME/etc/cron/r_hourly_runner.sh

# Run minute scripts (**)
* * * * * $APP_HOME/etc/cron/r_minute_runner.sh
